# localjudge/deb/omogenjudge-local.service
[Unit]
Description=OmogenJudge Local Judge
After=omogenjudge-sandbox.service
Requires=omogenjudge-sandbox.service

[Service]
Type=simple
Restart=always
User=omogenjudge-local
ExecStart=/usr/bin/omogenjudge-local

[Install]
WantedBy=multi-user.target
